# debian/kzorp-utils.cron.d
*/5 * * * *  root  [ -x /usr/sbin/kzorp-stats-updater ] && /usr/sbin/kzorp-stats-updater
